# .docker-vds-environment/images/fpm_npm_composer/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
directory=/app
command=...start